# src/modules/disable_usb/udev/10-disable_usb_storage.rules
ACTION=="add", KERNEL=="sd*", SUBSYSTEM=="block", ATTR{removable}=="1", RUN+="/usr/bin/eject %k"
